# bin/pironman5.service
# https://www.freedesktop.org/software/systemd/man/systemd.service.html
[Unit]
Description=pironman5 service
# Need to start last to avoid gpio being occupied
After=multi-user.target

[Service]
Type=forking
# WorkingDirectory=/opt/pironman5
ExecStart=/usr/local/bin/pironman5 start --background --debug-level=debug
ExecStop=/usr/local/bin/pironman5 stop
# PrivateTmp=False

[Install]
WantedBy=multi-user.target
